# infra/systemd/nous-backend.service
[Unit]
Description=nous-center backend (FastAPI + uvicorn)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=heygo
Group=heygo
WorkingDirectory=/media/heygo/Program/projects-code/_playground/nous-center/backend
EnvironmentFile=/media/heygo/Program/projects-code/_playground/nous-center/backend/.env
# uv lives in $HOME/.local/bin; systemd doesn't read .profile, point at it
Environment=PATH=/home/linuxbrew/.linuxbrew/bin:/home/heygo/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/linuxbrew/.linuxbrew/bin/uv run uvicorn src.api.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
# Log to journald — view with `journalctl -u nous-backend -f`
StandardOutput=journal
StandardError=journal
# Hard memory cap so a vLLM bug or runaway loader can't OOM the host.
# 8G (the pre-V1 value) was below wikeeyang fp8mixed's ~18G dequant peak —
# image_generate would get SIGKILL mid-load with no useful journal output.
# 64G leaves ~30G headroom on a 96G host for OS + cloudflared + ad-hoc
# tools. Bump to MemoryHigh=80G + MemoryMax=88G if you have more RAM and
# routinely load multiple image models concurrently.
MemoryMax=64G

[Install]
WantedBy=multi-user.target
