# infra/systemd/nous-status.service
[Unit]
Description=nous-center 独立状态监控(公开,纯 stdlib,2026-06-17)
# 刻意不 Requires/PartOf nous-backend —— 监控要在后端挂了时仍然活着显示「后端中断」。
After=network-online.target
Wants=network-online.target
# 监控自身崩溃环也别永久放弃。
StartLimitIntervalSec=0

[Service]
Type=simple
User=heygo
Group=heygo
# 系统 python3 即可(纯 stdlib,不碰 backend venv/torch)—— 真正解耦。
ExecStart=/usr/bin/python3 /media/heygo/Program/projects-code/_playground/nous-prod/infra/monitoring/status_service.py
Environment=NOUS_STATUS_HOST=127.0.0.1
Environment=NOUS_STATUS_PORT=8001
Environment=NOUS_STATUS_BACKEND=http://127.0.0.1:8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
