# infra/systemd/nous-healthprobe.timer
[Unit]
Description=Т»Ј 2 тѕєжњЪУиЉСИђТгА nous-center тЂЦт║итиАТБђ

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
AccuracySec=15s
Persistent=false

[Install]
WantedBy=timers.target
