# infra/systemd/nous-healthprobe.service
[Unit]
Description=nous-center 本地健康巡检(后端 + 公网隧道存活,2026-06-16)
# 一次性探针;由 nous-healthprobe.timer 触发。失败(后端连不上/DB 挂/隧道 down)
# 退出非 0 → systemd 标 failed;将来接告警只需加 OnFailure=<alert>.service。
After=network-online.target

[Service]
Type=oneshot
User=heygo
Group=heygo
ExecStart=/media/heygo/Program/projects-code/_playground/nous-center/infra/monitoring/nous-healthprobe.sh
StandardOutput=journal
StandardError=journal
